# roles/service/files/template.service
[Unit]
Description = %service_description%

[Service]
ExecStart = %service_start%
